# units/garden-systemd.service
[Unit]
Description=Garden Systemd Backend

[Service]
ExecStart=/usr/bin/garden-systemd --skeleton /var/lib/garden-systemd/skeleton
Restart=always

[Install]
WantedBy=multi-user.target
